# server/hermes-plugin/obsidian-sync/systemd/obsidian-hermes-queue.timer
[Unit]
Description=Run Obsidian Sync Hermes queue processor every minute

[Timer]
OnBootSec=60
OnUnitActiveSec=60
Unit=obsidian-hermes-queue.service

[Install]
WantedBy=timers.target
