# server/hermes-plugin/obsidian-sync/systemd/obsidian-hermes-queue.service
[Unit]
Description=Process one Obsidian Sync Hermes queue item
After=network-online.target

[Service]
Type=oneshot
User=hermes
WorkingDirectory=/home/obsidian-server
Environment=HOME=/home/hermes
Environment=HERMES_HOME=/home/hermes/.hermes
ExecStart=/home/hermes/.hermes/plugins/obsidian-sync/bin/process-queue.sh
